# scripts/supervisord/supervisor.conf
[program:run]
command = /var/lib/supervisord/bin/run
stdout_logfile=/dev/stdout
stdout_events_enabled=true

[inet_http_server]
port=localhost:9001
